# dataport-user-group-sync/debian/dataport-user-group-sync-sink.cron.d
PATH=/bin:/sbin:/usr/bin:/usr/sbin
*/5 * * * * root user_group_sync_sink.py
