# univention-ucc-status/ucc-send-status.conf
# UCC - Upstart script which starts univention-ucc-send-status
description     "Sending system information"

start on stopped ucc-syspol

task

exec    /usr/sbin/univention-ucc-send-status
